# root/etc/s6-overlay/s6-rc.d/app/run
#!/usr/bin/with-contenv bash
exec node server/Server.js